# goclaw/deploy/mango.service
[Unit]
Description=Mango Agent Gateway
After=network.target
 
[Service]
ExecStart=/usr/local/bin/mango serve
Restart=on-failure
RestartSec=5s
User=mango
Group=mango
RuntimeDirectory=mango
StateDirectory=mango
StandardOutput=journal
StandardError=journal
 
[Install]
WantedBy=multi-user.target

